# build/supervisord-syndic.conf
[supervisord]
nodaemon=true

[program:salt-master]
command=/usr/bin/salt-master
autostart=false
autorestart=true

[program:salt-syndic]
command=/usr/bin/salt-syndic
autostart=false
autorestart=true

[program:salt-minion]
command=/usr/bin/salt-minion
autostart=false
autorestart=true