# build/supervisord-syndic.conf
[supervisord]
nodaemon=true

[program:salt-master]
command=/usr/bin/salt-master -l debug
autostart=true
autorestart=true

[program:salt-syndic]
command=/usr/bin/salt-syndic -l debug
autostart=true
autorestart=true

[program:salt-minion]
command=/usr/bin/salt-minion -l debug
autostart=true
autorestart=true